# spec/sample_project/dockerfiles/mydockerrake_sshd/v2/supervisord/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
autorestart=true
